# deploy/supervisor/preshift-reverb.conf
; deploy/supervisor/preshift-reverb.conf
;
; Supervisor program definition for the Laravel Reverb WebSocket server.
; Reverb listens on port 6001 and is proxied through Nginx at /app (wss://).

[program:preshift-reverb]
command=php /var/www/preshift/api/artisan reverb:start --host=0.0.0.0 --port=6001
user=preshift
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
startsecs=5
stderr_logfile=/var/log/preshift/reverb.err.log
stdout_logfile=/var/log/preshift/reverb.out.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
