# deploy/supervisor/preshift-queue.conf
; deploy/supervisor/preshift-queue.conf
;
; Supervisor program definition for the Laravel queue worker.
; Processes queued jobs (notifications, broadcasts, etc.) with automatic restart.

[program:preshift-queue]
command=php /var/www/preshift/api/artisan queue:work --sleep=3 --tries=3 --max-time=3600
user=preshift
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
startsecs=5
stopwaitsecs=60
stderr_logfile=/var/log/preshift/queue.err.log
stdout_logfile=/var/log/preshift/queue.out.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
